# 附录/附录C/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
;UNIX socket文件路径，这里采用socket的方式管理而非端口，较之更为安全

[supervisord]
logfile=/tmp/supervisord.log  ;Supervisord日志文件路径
logfile_maxbytes=50MB  ;Supervisord日志文件大小，超出会rotate
logfile_backups=10  ;日志文件保留备份数量，默认值为10
loglevel=info             ;Supervisord日志级别，这里为info
pidfile=/tmp/supervisord.pid  ;Supervisord的pid文件路径
nodaemon=false            ;Supervisord以daemon 的方式运行
minfds=1024                ;可以打开的文件描述符的最小值，这里为1024
minprocs=200               ;可以打开的进程数的最小值，这里为200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface ;此处配置成默认值

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
;通过UNIX socket连接Supervisord
[program:index]
command=/usr/bin/python index.py  ;程序的启动命令
directory=/home/yhc/ContentEngine/api ;程序的启动目录
autostart=true ;在Supervisord启动的时候也自动启动
autorestart=true ;程序异常退出时自动重启
user=yhc ；程序用哪个用户启动
redirect_stderr=true ;把stderr重定向到stdout，即把错误日志重定向
stdout_logfile=/data/log/service_index.log ;sdout日志输出路径

[program:masterManager]
command=/usr/bin/python masterManager.py
directory=/home/yhc/ContentEngine
autostart=true
autorestart=true
user=yhc
redirect_stderr=true
stdout_logfile=/data/log/service_master.log
